# opsworks_resque_scheduler/templates/default/resque-scheduler.conf.erb
#!upstart
description "Resque scheduler"
respawn

env VERBOSE=true
export VERBOSE
env RACK_ENV=<%= node['opsworks_resque_scheduler']['rails_env'] %>
export RACK_ENV
env RAILS_ENV=<%= node['opsworks_resque_scheduler']['rails_env'] %>
export RAILS_ENV
env HOME=/home/deploy
export HOME

chdir <%= node['opsworks_resque_scheduler']['path'] %>
exec sudo -u deploy -E -- bundle exec rake environment --trace resque:scheduler PIDFILE=<%= node['opsworks_resque_scheduler']['pidfile'] %> >> log/resque-scheduler.log 2>&1

start on runlevel [2345]
stop on runlevel [016]
